# resources/autostart_on_startup.service
[Unit]
Description=Autostart bird service

[Service]
WorkingDirectory=/home/lars_rauer/bird
ExecStart=/home/lars_rauer/bird/start.sh
User=lars_rauer

[Install]
WantedBy=multi-user.target